# ebpf-program/my-bpf-sensor/debian/prerm
#!/bin/sh
# prerm script for my-bpf-sensor
#
# Stops the service before package removal.

set -e

case "$1" in
    remove|upgrade|deconfigure)
        if command -v systemctl >/dev/null 2>&1; then
            deb-systemd-helper stop my-bpf-sensor.service >/dev/null 2>&1 || true
        fi
        ;;
    failed-upgrade)
        ;;
    *)
        echo "prerm called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0